# dev/docker/supervisor/conf.d/jazzchords.conf
[program:jazzchords]
command=/usr/local/bin/uwsgi --ini=/etc/uwsgi/jazzchords.ini

